# deploy/todo-sync-update.timer
[Unit]
Description=Check GitHub for Todo Sync updates every minute

[Timer]
OnBootSec=30s
OnUnitActiveSec=60s
Unit=todo-sync-update.service

[Install]
WantedBy=timers.target
